# unraid old/UNRAID/config/plugins/dynamix/balance_cache.cron
# Generated btrfs balance cache schedule:
0 0 1 * * /usr/local/emhttp/plugins/dynamix/scripts/btrfs_balance start /mnt/cache -dusage=50 &> /dev/null

