# unraid old/UNRAID/config/plugins/dynamix/docker-update.cron
# Generated docker monitoring schedule:
10 0 * * 1 /usr/local/emhttp/plugins/dynamix.docker.manager/scripts/dockerupdate check &> /dev/null

