# unraid old/UNRAID/config/plugins/dynamix/mover.cron
# Generated mover schedule:
0 */1 * * * /usr/local/sbin/mover start |& logger -t move

