# unraid old/UNRAID/config/plugins/community.applications/notification_scan.cron

# CRON for CA background scanning of applications
23 * * * * php /usr/local/emhttp/plugins/community.applications/scripts/notices.php > /dev/null 2>&1

